# service/launchd/com.proj.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.proj.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.local/bin/proj</string>
        <string>daemon</string>
        <string>run</string>
    </array>

    <!-- The daemon launches `claude` (and resolves `proj`) when recreating
         pinned/keep-alive sessions; launchd gives only a minimal PATH, so add
         the user and Homebrew bin dirs where they are typically installed. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/.local/state/proj/daemon.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.local/state/proj/daemon.log</string>
</dict>
</plist>
