# service/systemd/proj-daemon.service
[Unit]
Description=proj-daemon; auto-resume Claude Code sessions after usage-limit cooldown
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/proj daemon run
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
